# srocket.service
[Unit]
Description=Sounding Rocket Mission Simulator
Documentation=https://github.com/kzaremski
After=network.target

[Service]
WorkingDirectory=/home/pi/sounding-rocket-tester/
Type=simple
User=pi
ExecStart=/usr/bin/python3 /home/pi/sounding-rocket-tester/app.js
Restart=on-failure

[Install]
WantedBy=multi-user.target